# docker-compose/supervisord_wechat.conf

[supervisord]
nodaemon=true

[program:wechat]
command=./wechat
directory=/app
autostart=true
autorestart=true
user=root
stdout_logfile=/app/logs/wechat.log


